# hyprsunset/.config/systemd/user/auto_hyprsunset.service
[Unit]
Description=Apply blue light filter

[Service]
Type=oneshot
ExecStart=/home/%u/.local/share/auto_hyprsunset
